# launchd-plists/com.personalscraper.backfill-ids.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
    com.personalscraper.backfill-ids.plist
    Weekly cross-provider IDs + ratings backfill (SH-3 / BD-S / DEV #28).

    Invokes `personalscraper library-backfill-ids` every Sunday at 03:00.
    The command walks media_item rows, fetches missing TMDB/TVDB/IMDb (via
    OMDb) / Rotten Tomatoes IDs and ratings, and merges them additively —
    never overwriting the canonical provider anchor or existing values.

    Note: `personalscraper library-backfill-ids` has no `--budget-seconds`
    flag (the original tech-debt 0.16.0 plan §8.1 mentioned one that does
    not exist on the actual CLI). The command is naturally bounded by the
    number of items remaining to backfill and TMDB/TVDB/OMDb rate limits;
    a full pass over a fresh library completes within a few minutes.

    Installation (per-user agent):
        # Ensure log directory exists (launchd does NOT create it):
        mkdir -p ~/.cache/personalscraper
        cp launchd-plists/com.personalscraper.backfill-ids.plist ~/Library/LaunchAgents/
        # Edit the copy to replace /Users/REPLACE_ME with your actual home path,
        # then load it:
        launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.personalscraper.backfill-ids.plist

    Removal:
        launchctl bootout gui/$(id -u) ~/Library/LaunchAgents/com.personalscraper.backfill-ids.plist
        rm ~/Library/LaunchAgents/com.personalscraper.backfill-ids.plist

    Prerequisites (verified by the command itself; agent will exit non-zero
    if missing):
      - `personalscraper library-init-canonical` must have run at least once
        so canonical_provider is seeded on legacy rows.
      - `.env` must contain `TMDB_API_KEY`, `TVDB_API_KEY`, and `OMDB_API_KEY`.

    Trigger a manual run for testing:
        launchctl kickstart -k gui/$(id -u)/com.personalscraper.backfill-ids
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.personalscraper.backfill-ids</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Replace /usr/local/bin with the output of: which personalscraper -->
        <string>/usr/local/bin/personalscraper</string>
        <string>library-backfill-ids</string>
    </array>

    <!-- Logs land in ~/.cache/personalscraper/ — the directory MUST exist
         before launchd fires (mkdir -p step in Installation above). -->
    <key>StandardOutPath</key>
    <string>/Users/REPLACE_ME/.cache/personalscraper/backfill-ids.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/REPLACE_ME/.cache/personalscraper/backfill-ids.log</string>

    <!-- Weekly: every Sunday at 03:00. Weekday: 0 == Sunday in launchd.
         03:00 puts it after the nightly index passes (03:30 quick, etc.). -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Do NOT run immediately on agent load; only fire at the scheduled time. -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Working directory must contain config/ so the CLI can locate
         config.json5 + .env. -->
    <key>WorkingDirectory</key>
    <string>/Users/REPLACE_ME/dev/PersonnalScaper</string>

    <!-- Inherit a sane PATH so pyenv/Homebrew shims resolve. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
