# docs/reference/launchd/personalscraper-index-enrich.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
    personalscraper-index-enrich.plist
    Weekly enrich pass on Sunday at 04:00.
    Runs pymediainfo + NFO parse + artwork inventory on files that have never
    been enriched (or whose enrichment is stale).  Budget is capped at 1800 s
    (30 min) so the job finishes before the morning pipeline cron.

    This is an OPTIONAL plist.  Skip it if you run the enrich mode manually
    or if your disks are fast enough that enrich is included in the
    nightly full-disk rotate.

    Installation (per-user agent):
        cp personalscraper-index-enrich.plist ~/Library/LaunchAgents/
        launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/personalscraper-index-enrich.plist

    Removal:
        launchctl bootout gui/$(id -u) ~/Library/LaunchAgents/personalscraper-index-enrich.plist
        rm ~/Library/LaunchAgents/personalscraper-index-enrich.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.personalscraper.index-enrich</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Replace /usr/local/bin with the output of: which personalscraper -->
        <string>/usr/local/bin/personalscraper</string>
        <string>library-index</string>
        <string>--mode</string>
        <string>enrich</string>
        <string>--budget</string>
        <string>1800</string>
        <string>--wait-for-lock</string>
        <string>0</string>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/REPLACE_ME/dev/PersonnalScaper/__logit__/index.enrich.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/REPLACE_ME/dev/PersonnalScaper/__logit__/index.enrich.err.log</string>

    <!-- Weekly on Sunday at 04:00, well before the morning pipeline cron. -->
    <key>StartCalendarInterval</key>
    <dict>
        <!-- Weekday 0 = Sunday (launchd convention). -->
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>LaunchOnlyOnce</key>
    <false/>

    <!-- Working directory must contain config.json5. -->
    <key>WorkingDirectory</key>
    <string>/Users/REPLACE_ME/dev/PersonnalScaper</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
