# docs/reference/launchd/personalscraper-index-rotate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
    personalscraper-index-rotate.plist
    Nightly rotating full-disk index at 03:35.
    Delegates to index-rotate.sh which selects the disk by weekday:
        Mon -> Disk1, Tue -> Disk2, Wed -> Disk3, Thu -> Disk4
        Fri/Sat/Sun -> falls back to quick mode (no disk argument)

    This plist runs IN ADDITION to the quick plist, not instead of it.
    Use when you want a guaranteed full-rebuild cadence (one disk per night).

    Installation (per-user agent):
        cp personalscraper-index-rotate.plist ~/Library/LaunchAgents/
        launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/personalscraper-index-rotate.plist

    Removal:
        launchctl bootout gui/$(id -u) ~/Library/LaunchAgents/personalscraper-index-rotate.plist
        rm ~/Library/LaunchAgents/personalscraper-index-rotate.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.personalscraper.index-rotate</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Replace the path with the actual location of index-rotate.sh.
             Make the script executable first: chmod +x index-rotate.sh -->
        <string>/bin/bash</string>
        <string>/Users/REPLACE_ME/dev/PersonnalScaper/docs/reference/launchd/index-rotate.sh</string>
    </array>

    <!-- Logs for the wrapper script itself (disk selection logic).
         personalscraper output goes to its own log inside __logit__/. -->
    <key>StandardOutPath</key>
    <string>/Users/REPLACE_ME/dev/PersonnalScaper/__logit__/index.rotate.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/REPLACE_ME/dev/PersonnalScaper/__logit__/index.rotate.err.log</string>

    <!-- Run at 03:35 every night, 5 minutes after the quick scan. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>35</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>LaunchOnlyOnce</key>
    <false/>

    <!-- Working directory must contain config.json5. -->
    <key>WorkingDirectory</key>
    <string>/Users/REPLACE_ME/dev/PersonnalScaper</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
